# agent-business-os/ops/launchd/com.famtastic.abos-tick.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.famtastic.abos-tick</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/env</string>
    <string>node</string>
    <string>{{HUB}}/agent-business-os/agents/run.js</string>
    <string>tick</string>
  </array>
  <key>WorkingDirectory</key><string>{{HUB}}/agent-business-os</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>ABOS_LIVE_BILLING</key><string>1</string>
    <key>ABOS_AUTO_CLOSE</key><string>1</string>
    <key>PATH</key><string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
  </dict>
  <key>StartInterval</key><integer>900</integer>
  <key>RunAtLoad</key><true/>
  <key>StandardOutPath</key><string>/tmp/abos-tick.log</string>
  <key>StandardErrorPath</key><string>/tmp/abos-tick.err</string>
</dict>
</plist>
